# root/etc/uci-defaults/luci-theme-solusi
#!/bin/sh
	uci batch <<-EOF
        set luci.themes.Solusi=/luci-static/solusi
        commit luci
	EOF
exit 0